# backend/deploy/systemd/vigia-bootstrap.service
# Exemplo de unit systemd para vigia-bootstrap (Linux).
# Podes gerar e instalar automaticamente (recomendado), com o binário no PATH ou com caminho completo:
#   sudo /caminho/para/vigia-bootstrap install-service -data-dir /var/lib/vigia/bootstrap -now
#
# Instalação manual do ficheiro:
#   sudo install -m644 vigia-bootstrap.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vigia-bootstrap.service

[Unit]
Description=Vigia bootstrap (OTA / Docker compose)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
# Garante que o data-dir existe antes do arranque
ExecStartPre=/bin/mkdir -p /var/lib/vigia/bootstrap
# Estado persistente: compose, bootstrap.yaml, device-id
ExecStart=/usr/local/bin/vigia-bootstrap -data-dir /var/lib/vigia/bootstrap
# Opcional: variáveis sem editar o unit (criar o ficheiro em /etc/vigia/bootstrap.env)
# EnvironmentFile=-/etc/vigia/bootstrap.env
# Ex.: VIGIA_DEVICE_ID=..., VIGIA_AUTO_INSTALL=0 após Docker já instalado na imagem
Restart=on-failure
RestartSec=5
# Primeira execução com auto-install pode precisar de root para apt/get.docker.com;
# depois pode usar um utilizador no grupo docker e ajustar permissões em /var/lib/vigia.
User=root

[Install]
WantedBy=multi-user.target
